# provision/conf/supervisor/programs/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
stdout_logfile=/opt/app/logs/nginx/access.log
stderr_logfile=/opt/app/logs/nginx/error.log
